# home/.config/systemd/user/restic-mount@.service
[Unit]
Description=Restic mount service

[Service]
Type=exec
ExecStart=/usr/bin/restic mount %h/%I

